# tg_my_home_scan_bot_watcher.path
[Path]
PathModified=/srv/http/TGMyHomeScanBot/src/.git

[Install]
WantedBy=multi-user.target
